# backend/app/supervisor/supervisord.d/laravel-fpm.conf
[program:laravel-fpm]
process_name=laravel-fpm
command=php-fpm
autostart=true
autorestart=true
numprocs=1
user=root
redirect_stderr=true
stdout_logfile=/var/www/server/storage/logs/php-fpm.log